# sixthsense/platform/linux/shipper/service/log/run
#!/bin/bash
exec multilog t /usr/local/logmind/sixthsense/logs/shipper/service
