# sixthsense/platform/linux/indexer/service/log/run
#!/bin/bash
exec multilog t /usr/local/logmind/sixthsense/logs/indexer/service
