# openmind2/platform/linux/service/log/run
#!/bin/bash
exec multilog t /usr/local/logmind/openmind/logs/service
