# elasticsearch/platform/linux/service/log/run
#!/bin/bash
exec multilog t /usr/local/logmind/elasticsearch/logs/service
